# skills/newapi-usage/assets/com.newapi.usage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.newapi.usage</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__SKILL_DIR__/scripts/report.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/newapi-usage.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/newapi-usage.err.log</string>
</dict>
</plist>
